# root/etc/systemd/system/iwd.service.d/iwd-delay.conf
[Service]
ExecStartPre=/usr/bin/sleep 3
